# sway/gammastep.service
[Unit]
Description=Display colour temperature adjustment
PartOf=graphical-session.target

[Service]
ExecStart=/usr/bin/gammastep
Restart=on-failure

[Install]
WantedBy=sway-session.target
