# config/systemd/user/clean-nvim-backups.timer
[Unit]
Description=Timer for cleaning up neovim's backup and swap files

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
