# config/systemd/user/firefox-profile-on-ram.service
[Unit]
Description=Firefox profile in memory

[Install]
WantedBy=default.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%E/systemd/user/firefox-sync.sh
ExecStop=%E/systemd/user/firefox-sync.sh
