# sway/low-battery/low-battery.service
[Unit]
Description=low battery notification service for sway using swaynag

[Service]
Type=forking
ExecStart=%E/systemd/user/low-battery.sh
